# _server_setup/periodic_archiveooo.crontab
# m h  dom mon dow   command
*/5 *    * * *   ~/archiveooo/env/bin/python3 -u ~/archiveooo/delete_aged_vms.py --log-level=WARNING
